# deploy/dolphin.service
[Unit]
Description=Dolphin AI Agent
Documentation=https://github.com/dolphinZzv/dolphin
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dolphin
Group=dolphin

# Binary and config
ExecStart=/usr/local/bin/dolphin
Restart=on-failure
RestartSec=5

# Working directory (session logs, skills, crontab stored here)
WorkingDirectory=/var/lib/dolphin

# Environment variables for LLM config
Environment=DZ_LLM_API_KEY=
Environment=DZ_LLM_MODEL=
Environment=DZ_LLM_BASE_URL=

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ProtectKernelTunables=true
ProtectControlGroups=true
ReadWritePaths=/var/lib/dolphin

# File descriptor limits
LimitNOFILE=65536

# Logging (journald)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
